# templates/restic.service.j2
[Unit]
Description=Backup {{ item.name }} using restic

[Service]
Type=oneshot
ExecStart={{ restic_script_dir }}/backup-{{ item.name }}.sh
TimeoutStartSec=0
Environment="CRON=true"
